# services/vaultwarden/quadlet/vaultwarden.container.tmpl
[Unit]
Requires=vaultwarden-ts.service
After=vaultwarden-ts.service

[Container]
Image=docker.io/vaultwarden/server:latest
AutoUpdate=registry
Volume=vaultwarden-data.volume:/data
Network=container:systemd-vaultwarden-ts
Environment=SIGNUPS_ALLOWED=true
Environment=DOMAIN=https://vaultwarden.{{ .Env.TS_DOMAIN }}

[Service]
Restart=always

[Install]
WantedBy=default.target
