# services/ollama/systemd/ollama-ts-watchdog.timer
[Unit]
Description=Check ollama-ts health periodically

[Timer]
OnBootSec=2min
OnUnitActiveSec=2min

[Install]
WantedBy=timers.target
